# programs/web-playground/deploy/playground.service
[Unit]
Description=Objeck Playground API
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
User=playground
Group=playground
WorkingDirectory=/opt/playground/backend
ExecStart=/opt/playground/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

Environment=PLAYGROUND_SANDBOX_IMAGE=objeck-sandbox:latest
Environment=PLAYGROUND_HOST_TMP_DIR=/tmp/playground
Environment=PLAYGROUND_SHARE_DB_PATH=/opt/playground/data/shares.db
Environment=PLAYGROUND_ALLOWED_ORIGINS=["https://www.objeck.org","https://objeck.org","https://playground.objeck.org"]

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/playground/data /tmp/playground

[Install]
WantedBy=multi-user.target
